# iptables-port-8080.service
[Unit]
Description=iptables rules for port 8080 access control
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/home/exx/myCode/genai-agent-core/setup-port-8080-firewall.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
